# .config/systemd/pacman-fetch.timer
[Unit]
Description=Fetch upstream package infos

[Timer]
OnCalendar=00/1:00

[Install]
WantedBy=timers.target
