# .config/systemd/pacman-fetch.service
[Unit]
Description=Fetch upstream package infos

[Service]
Type=oneshot
ExecStart=pacman -Sy

[Install]
WantedBy=
